# deploy/kop-whatsapp-bridge.service
[Unit]
Description=KOP-ML WhatsApp Bridge
After=network.target

[Service]
Type=simple
User=kop
Group=kop
WorkingDirectory=/opt/kop-ml/whatsapp_bridge
Environment=PUPPETEER_EXECUTABLE_PATH=/usr/bin/chromium-browser
Environment=PUPPETEER_SKIP_CHROMIUM_DOWNLOAD=true
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
