# deploy/kop-scraper.service
[Unit]
Description=KOP-ML Scraper (Pelando + Mercado Livre)
After=network.target kop-whatsapp-bridge.service
Wants=kop-whatsapp-bridge.service

[Service]
Type=simple
User=kop
Group=kop
WorkingDirectory=/opt/kop-ml
EnvironmentFile=/opt/kop-ml/.env.prod
Environment=APP_ENV=prod
Environment=HEADLESS=true
# CHROME_BINARY vazio: scraper/browser.py usa chromium do Playwright
# (~/.cache/ms-playwright/chromium-*/chrome-linux/chrome). O snap chromium do
# Ubuntu é incompatível com nodriver (confinamento bloqueia CDP).
Environment=CHROME_BINARY=
ExecStart=/opt/kop-ml/venv/bin/python main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
